# 00-local/opencde/profile.d/opencde.csh
#!/bin/csh
setenv PATH ${PATH}:/usr/dt/bin
