# 00-alternate/openmotif/profile.d/openmotif.csh
#!/bin/csh
setenv PATH ${PATH}:/opt/openmotif/bin
setenv MANPATH ${MANPATH}:/opt/openmotif/man
