# 00-alternate/openmotif/profile.d/openmotif.sh
#!/bin/sh
export PATH=${PATH}:/opt/openmotif/bin
export MANPATH=${MANPATH}:/opt/openmotif/man
