# 00-alternate/openmotif/profile.d/openmotif.csh
#!/bin/csh
setenv PATH ${PATH}:/usr/dt/bin
setenv MANPATH ${MANPATH}:/usr/dt/man
